# cookbooks/chub_yui_combo/files/default/combo.conf
setuid combo
setgid combo
exec combohandler --port 8000 --rootsFile /etc/combo/routes.json --webRoot "/var/combo"
respawn
start on runlevel [2345]
stop on runlevel [016]
