# cookbooks/yui_combo/files/default/combo.conf
setuid combo
setgid combo
exec combohandler -f /var/combo/routes.json -w "/var/combo"
respawn
start on runlevel [2345]
stop on runlevel [016]