# cookbooks/chub-vault/templates/default/vault.conf.erb
setuid vault
setgid vault
chdir /var/vault
exec java -server -d64 -Xms<%= @java_heap_size %> -Xmx<%= @java_heap_size %> -Djava.net.preferIPv4Stack=true -jar /opt/vault/vault.jar server /etc/vault/vault.yaml
respawn
start on runlevel [2345]
stop on runlevel [016]
