# cookbooks/chub_mc_webservices/templates/default/connman_webservice.conf.erb
exec java -jar <%= node[:chub_mc_webservice][:connman][:deploy_dir] %>/<%= node[:chub_mc_webservice][:connman][:jar_file_name] %> server <%= node[:chub_mc_webservice][:connman][:deploy_dir] %>/<%= node[:chub_mc_webservice][:connman][:config_file_name] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
