# cookbooks/chub_mc_webservices/templates/default/flight_webservice.conf.erb
exec java -jar <%= node[:chub_mc_webservice][:flight][:deploy_dir] %>/<%= node[:chub_mc_webservice][:flight][:jar_file_name] %> server 
respawn
start on runlevel [2345]
stop on runlevel [016]
