# cookbooks/chub_envdash/templates/default/envdash.conf.erb
setuid envdash
setgid envdash
chdir /var/envdash
script
  HOSTNAME=$(hostname -s)
  exec java -server -d64 -Xms<%= node[:chub_envdash][:java_heap_size] %> -Xmx<%= node[:chub_envdash][:java_heap_size] %> -Dratpack.other.title="<%= node[:chub_envdash][:title] %>" -Dratpack.other.prefix="<%= node[:chub_envdash][:prefix] %>" -Dratpack.other.environmentNames="<%= node[:chub_envdash][:environment_names].join(' ') %>" -Dratpack.other.chefClient=$HOSTNAME -Dratpack.other.chefPem=/etc/chef/client.pem -Dratpack.other.sensuUsername="<%= node[:sensu][:api][:user] %>" -Dratpack.other.sensuPassword="<%= node[:sensu][:api][:password] %>" -Dratpack.port=8080 -jar /opt/envdash/envdash.jar
end script
respawn
start on runlevel [2345]
stop on runlevel [016]
