# cookbooks/chub_envdash/templates/default/envdash.conf.erb
setuid envdash
setgid envdash
chdir /var/envdash
exec java -server -d64 -Xms<%= node[:chub_envdash][:java_heap_size] %> -Xmx<%= node[:chub_envdash][:java_heap_size] %> -Dratpack.other.title="<%= node[:chub_envdash][:title] %>" -Dratpack.other.prefix="<%= node[:chub_envdash][:prefix] %>" -Dratpack.other.environmentNames="<%= node[:chub_envdash][:environment_names].join(' ') %>" -jar /opt/envdash/envdash.jar
respawn
start on runlevel [2345]
stop on runlevel [016]
