# cookbooks/chub_landingstrip/templates/default/landingstrip.conf.erb
exec java -jar <%= node[:chub_landingstrip][:app][:deploy_dir] %>/<%= node[:chub_landingstrip][:app][:app_name] %>.jar context=<%= node[:chub_landingstrip][:app][:app_name] %> port=<%= node[:chub_landingstrip][:app][:port] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
