# cookbooks/chub_browserstack_tunnel/templates/default/browserstack-tunnel.conf.erb
setuid browserstack
setgid browserstack
chdir /var/browserstack-tunnel
exec /opt/browserstack-tunnel/BrowserStackLocal <%= node[:chub_browserstack_tunnel][:api_key] %> <% node[:chub_browserstack_tunnel][:services].to_enum.with_index(1) do |service, index| %><%= service.host %>,<%= service.port %>,<%= service.ssl ? 1 : 0 %><% if(index < node[:chub_browserstack_tunnel][:services].length) %>,<% end %><% end %> -skipCheck
respawn
start on runlevel [2345]
stop on runlevel [016]
