# cookbooks/chub_browserstack_tunnel/templates/default/browserstack-tunnel.conf.erb
setuid browserstack
setgid browserstack
chdir /var/browserstack-tunnel
exec /opt/browserstack-tunnel/BrowserStackLocal <%= node[:chub_browserstack_tunnel][:api_key] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
