# cookbooks/chub-browserstack-tunnel/templates/default/browserstack-tunnel.conf.erb
setuid browserstack
setgid browserstack
chdir /var/browserstack-tunnel
exec java -jar /opt/browserstack-tunnel/BrowserStackTunnel.jar <%= node["chub-browserstack-tunnel"]["api_key"] %> <% node['chub-browserstack-tunnel']['services'].each do |service| %><%= service.host %>,<%= service.port %>,<%= service.ssl ? 1 : 0 %>,<% end %>
respawn
start on runlevel [2345]
stop on runlevel [016]
