# cookbooks/chub-klerk/templates/default/klerk.conf.erb
setuid klerk
setgid klerk
chdir /var/opt/klerk
exec java -server -d64 -Xms<%= @java_heap_size %> -Xmx<%= @java_heap_size %> -jar /opt/klerk/klerk-app.jar server /etc/opt/klerk/klerk.yml
respawn
start on runlevel [2345]
stop on runlevel [016]
