# cookbooks/chub_mc_authservice/templates/default/mc_authservice.conf.erb
chdir <%= node[:chub_mc_authservice][:deploy_dir] %>
exec java -jar /<%= node[:chub_mc_authservice][:deploy_dir] %>/<%= node[:chub_mc_authservice][:jar_file_name] %>.jar -httpPort <%= node[:chub_mc_authservice][:httpPort] %> -httpsPort <%= node[:chub_mc_authservice][:httpsPort] %> -keyAlias <%= node[:chub_mc_authservice][:keyAlias] %> -Djavax.net.ssl.keyStore=../<%= node[:chub_mc_authservice][:keyStore] %> -Djavax.net.ssl.keyStorePassword=<%= node[:chub_mc_authservice][:keyStorePassword] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
