# cookbooks/chub_onboarding_webservice/templates/default/onboarding_webservice.conf
chdir <%= node["chub_onboarding_webservice"]["deployDir"] %>
exec java -jar <%= node["chub_onboarding_webservice"]["jarName"] %> -server <%= node["chub_onboarding_webservice"]["propertyFileName"] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
