# cookbooks/chub-onboarder/templates/default/onboarder.conf.erb
chdir <%= node['chub-onboarder']['app']['deploy_dir'] %>
exec java -jar /<%= node['chub-onboarder']['app']['deploy_dir'] %>/<%= node['chub-onboarder']['app']['app_name'] %>.jar context=<%= node['chub-onboarder']['app']['app_name'] %> port=<%= node['chub-onboarder']['app']['port'] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
