# cookbooks/chub_mc_app/templates/default/app.conf.erb
exec java<% node[:chub_mc_app][:config_vars].each do |key, value| -%><% unless value.empty? -%> -D<%= key %>=<%= value %><% end %><% end %> -jar <%= node[:chub_mc_app][:deploy_dir] %>/<%= node[:chub_mc_app][:app_name] %>.jar context=<%= node[:chub_mc_app][:app_name] %> port=<%= node[:chub_mc_app][:port] %>
respawn
start on runlevel [2345]
stop on runlevel [016]
