# cookbooks/chub_oddrest/templates/default/oddrest.conf.erb
setuid oddrest
setgid oddrest
chdir /var/oddrest
exec java -server -d64 -Xms<%= node[:chub_oddrest][:java_heap_size] %> -Xmx<%= node[:chub_oddrest][:java_heap_size] %> -Djava.net.preferIPv4Stack=true -Dcom.sun.management.jmxremote.port=<%= node[:chub_oddrest][:jmxremote_port] %> -Dcom.sun.management.jmxremote.authenticate=false -Dcom.sun.management.jmxremote.ssl=false -DorgDataDirOverride.config=/etc/oddrest/override.properties -DorgDataDir.config=<%= node[:chub_oddrest][:properties_location] %> -jar /opt/oddrest/oddrest.jar server /etc/oddrest/oddrest.yaml
respawn
start on runlevel [2345]
stop on runlevel [016]